# panel/files/etc/services.d/pterodactyl/run
#!/usr/bin/with-contenv bash

cd /var/www/pterodactyl || exit

exec s6-setuidgid abc nginx –g "daemon off;"