# panel/files/etc/services.d/phpfpm/run
#!/usr/bin/with-contenv bash

exec php-fpm7.2 -F